# honcho_export_systemd/templates/systemd/master.target
[Unit]
StopWhenUnneeded=true
Wants={{ master_wants }}

[Install]
WantedBy=multi-user.target
