# honcho_export_systemd/templates/systemd/process_master.target
[Unit]
StopWhenUnneeded=true
Wants={{ process_master_wants }}

